# etc/systemd/user/boost-control.service
[Unit]
Description=Toggle boost alongside power profile changes

[Service]
ExecStart=/usr/local/bin/boost-control

[Install]
WantedBy=default.target
